# base/bash/.config/profile.d/70_fzf.sh
#!/usr/bin/env bash

if [[ -r ~/.fzf.bash ]]; then
    source ~/.fzf.bash
fi
