# base/bash/.config/profile.d/999_dummy.sh
#!/usr/bin/env bash

# My shell prompt displays a frowny face when the last process exited with an
# error status.  Since some of the files in .dotfiles/source test for OS, I
# often get the frowny face when, e.g., spawning a new shell.  This dummy file
# silences the error.

return 0
